# board/pure/rootfs-overlay/etc/sv/squeeze2upnp/run
#!/bin/sh

exec  /usr/aprenderer/upnp/squeeze2upnp -z -k -x /usr/aprenderer/upnp/config.xml

